# etc/supervisord.conf.in
[unix_http_server]
file=${buildout:directory}/temp/supervisor.sock ; path to your socket file

[supervisorctl]
serverurl=unix://${buildout:directory}/temp/supervisor.sock ; use a unix:// URL for a unix socket

[supervisord]
logfile = ${buildout:directory}/log/supervisord.log
logfile_maxbytes = 50MB
logfile_backups = 10
loglevel = info
pidfile = ${buildout:directory}/temp/supervisord.pid
nodaemon = false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:oauth]
command = ${buildout:directory}/bin/oauth
autostart = false
autorestart = true
stderr_logfile = ${buildout:directory}/log/oauth-stderr.log
stdout_logfile = ${buildout:directory}/log/oauth-stdout.log

[program:api]
command = ${buildout:directory}/bin/gunicorn -w 2 -k tornado api.manage:application --log-level=info -b 127.0.0.1:5000
autostart = false
autorestart = true
stderr_logfile = ${buildout:directory}/log/api-stderr.log
stdout_logfile = ${buildout:directory}/log/api-stdout.log
loglevel = info

[program:redis_01]
autostart = false
autorestart = true
command = redis-server ${buildout:directory}/etc/redis_01.conf
stderr_logfile = ${buildout:directory}/log/redis_01-stderr.log
stdout_logfile = ${buildout:directory}/log/redis_01-stdout.log

[program:redis_02]
autostart = false
autorestart = true
command = redis-server ${buildout:directory}/etc/redis_02.conf
stderr_logfile = ${buildout:directory}/log/redis_02-stderr.log
stdout_logfile = ${buildout:directory}/log/redis_02-stdout.log

[program:demo]
autostart = false
autorestart = true
command = ${buildout:directory}/bin/demo
stderr_logfile = ${buildout:directory}/log/demo-stderr.log
stdout_logfile = ${buildout:directory}/log/demo-stdout.log